# deployments/systemd/messagebridge.service
[Unit]
# Описание сервиса
Description=MessageBridge - Secure Payment Webhook Gateway
# Ссылка на документацию
Documentation=https://github.com/expai/messagebridge
# Запускать после инициализации сети         
After=network.target
# Желательно, чтобы сеть была доступна                                        
Wants=network.target

[Service]
# Обычный сервис, основной процесс — это ExecStart
Type=simple
# Запускать от пользователя messagebridge
User=messagebridge
# Запускать от группы messagebridge
Group=messagebridge

# Security settings
# Запрещает процессу получать новые привилегии (например, через setuid)
NoNewPrivileges=true
# Выделяет отдельный /tmp для сервиса (невидим для других процессов)
PrivateTmp=true
# Только чтение для всей системы, кроме явно разрешённых путей      
ProtectSystem=strict
# Запрещает доступ к домашним директориям пользователей
ProtectHome=true
# Разрешает запись только в эти директории
ReadWritePaths=/var/lib/messagebridge /var/log/messagebridge
# Запрещает изменение настроек ядра (sysctl)
ProtectKernelTunables=true
# Запрещает загрузку/выгрузку модулей ядра
ProtectKernelModules=true
# Запрещает доступ к cgroups (контролю ресурсов)
ProtectControlGroups=true
# Запрещает использование real-time приоритетов (RT-процессов)
RestrictRealtime=true
# Запрещает использование файлов с SUID/SGID битами
RestrictSUIDSGID=true
# Удаляет IPC-ресурсы (shared memory, semaphores) при завершении процесса
RemoveIPC=true
# Скрывает /dev, только минимальный набор устройств доступен
PrivateDevices=true

# Environment
# Переменная окружения с путём к конфигу
Environment=CONFIG_PATH=/etc/messagebridge/config.yaml
# Рабочая директория процесса
WorkingDirectory=/var/lib/messagebridge

# Process management
# Команда запуска сервиса
ExecStart=/usr/local/bin/messagebridge -config /etc/messagebridge/config.yaml
# Как перезапустить (reload) сервис
ExecReload=/bin/kill -HUP $MAINPID
# Сначала SIGTERM основному процессу, потом дочерним
KillMode=mixed
# Сигнал для остановки процесса
KillSignal=SIGTERM
# Время ожидания завершения процесса при остановке (сек)
TimeoutStopSec=30

# Restart policy
# Всегда перезапускать при завершении
Restart=always
 # Ждать 10 секунд перед перезапуском
RestartSec=10
# За 300 секунд (5 минут)...
StartLimitInterval=300
# ...разрешено не более 5 рестартов (иначе сервис будет заблокирован)
StartLimitBurst=5

# Restart conditions
# Сигнал для завершения при рестарте
RestartKillSignal=SIGTERM
# Не рестартовать, если процесс завершился с кодом 0 (успех)
RestartPreventExitStatus=0

# Handle exit codes
# Коды завершения, считающиеся успешными (0 — обычный успех, 2 — сигнализирует о необходимости рестарта)
SuccessExitStatus=0 2
# Если приложение завершилось с кодом 2 — всегда рестартовать
RestartForceExitStatus=2

# Logging
# Логировать stdout в systemd journal
StandardOutput=journal
# Логировать stderr в systemd journal
StandardError=journal
# Имя сервиса в логах
SyslogIdentifier=messagebridge

# Resource limits
# Максимальное количество открытых файловых дескрипторов
LimitNOFILE=65536
# Максимальное количество процессов/потоков
LimitNPROC=4096

[Install]
# Автоматически запускать при загрузке системы (обычный режим для серверов)
WantedBy=multi-user.target 